# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]  ; 新增SSH服务
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
priority=1       ; 调整优先级确保最先启动
startsecs=5      ; 增加启动等待时间
stopsignal=TERM

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x800x24
autostart=true
autorestart=true
priority=2       ; 调整优先级
startretries=3   ; 增加重试次数

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -nopw -forever -shared -rfbport 5900
autostart=true
autorestart=true
priority=3
depends_on=xvfb  ; 显式声明依赖关系

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ 6080 localhost:5900
autostart=true
autorestart=true
priority=4
depends_on=x11vnc

[program:runit]
command=/usr/bin/runsvdir -P /etc/service
autostart=true
autorestart=true
priority=6

[program:elisim]
command=/bin/bash -c 'export DISPLAY=:1 && source /etc/profile && cd /home/elite && bash run_elisim.sh -robotTypeId="$ROBOT_TYPE_ID" -toolTypeId="$TOOL_TYPE_ID" -controlBoxId="$CONTROL_BOX_ID" -robotGenerationId="$ROBOT_GENERATION_ID"'
autostart=true
autorestart=true
priority=5
environment=ROBOT_TYPE_ID=%(ENV_ROBOT_TYPE_ID)s,TOOL_TYPE_ID=%(ENV_TOOL_TYPE_ID)s,CONTROL_BOX_ID=%(ENV_CONTROL_BOX_ID)s,RobotGenerationId=%(ENV_ROBOT_GENERATION_ID)s
depends_on=xvfb,x11vnc
